# prometheus-raspberry-exporter.service
[Unit]
Description=Collect Raspberry Pi metrics for Prometheus

[Service]
Type=oneshot
ExecStart=/usr/local/bin/raspberry_exporter
