# prometheus-raspberry-exporter.timer
[Unit]
Description=Temperature CPU Timer

[Timer]
OnCalendar=*:*:0,15,30,45

[Install]
WantedBy=timers.target
