# raspberry-metrics.service
[Unit]
Description=Collect Raspberry Pi metrics for Prometheus

[Service]
Type=oneshot
ExecStart=/home/teemow/raspberry-metrics.sh
